# perfsonar-checker-systemd/ps-checker-daemon.service
[Unit]
Description=perfSONAR-Checker

[Service]
ExecStart=/usr/local/bin/run-perfsonar-tests.sh

[Install]
WantedBy=multi-user.target
